# deploy/systemd/catty-app.service
[Unit]
Description=Catty reminders app
After=network.target

[Service]
Type=simple
User=lenya
Group=lenya
WorkingDirectory=/opt/catty-reminders-app
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/catty-reminders-app/.env
ExecStart=/opt/catty-reminders-app/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8181
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
